# services/hyperion-claude.service
[Unit]
Description=Hyperion Claude - Always-on Claude Code session
After=network.target hyperion-router.service
Wants=hyperion-router.service

[Service]
Type=forking
User=admin
Group=admin
WorkingDirectory=/home/admin/hyperion-workspace
Environment=PATH=/home/admin/.claude/bin:/home/admin/.local/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOME=/home/admin
ExecStart=/usr/bin/tmux -L hyperion new-session -d -s hyperion -c /home/admin/hyperion-workspace /home/admin/hyperion/scripts/claude-wrapper.exp
ExecStop=/usr/bin/tmux -L hyperion kill-session -t hyperion
RemainAfterExit=yes
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
