# services/lobster-observability.service
[Unit]
Description=Lobster Observability HTTP Server (port 8742)
After=network.target

[Service]
Type=simple
User=lobster
Group=lobster
WorkingDirectory=/home/lobster/lobster
Environment=PATH=/home/lobster/.local/bin:/home/lobster/.cargo/bin:/usr/local/bin:/usr/bin:/bin
Environment=LOBSTER_CONFIG_DIR=/home/lobster/lobster-config
Environment=LOBSTER_INSTALL_DIR=/home/lobster/lobster
Environment=LOBSTER_MESSAGES=/home/lobster/messages
Environment=LOBSTER_WORKSPACE=/home/lobster/lobster-workspace
EnvironmentFile=-/home/lobster/lobster/config/config.env
EnvironmentFile=/home/lobster/lobster-config/config.env
ExecStart=/home/lobster/lobster/.venv/bin/lobster-observability-server --port 8742
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=lobster-observability

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target